# tmux/dakdjob.conf
description "Dakd project job"
author "Pau"
start on runlevel [2345]
stop on shutdown

exec su -l user -s /bin/sh -c "/home/user/project/tmux/dashboard.sh"
