# tmux/dakd.service
[Unit]
Description=dakd malware analysis project
After=network.target

[Service]
Type=forking
User=homero
group=homero
ExecStart=/home/user/dakd-project/tmux/dashboard.sh
ExecStop=/usr/bin/tmux kill-session -t dakd

[Install]
WantedBy=multi-user.target
